# services/service_overrides/hermes-gateway.service.d/30-cost-guardrail.conf
[Service]
Environment=HERMES_COST_GUARDRAIL_ENABLED=1
Environment=HERMES_DAILY_COST_LIMIT_USD=2.00
Environment=HERMES_COST_POSTGRES_DB=rag
