# services/service_overrides/hermes-gateway.service.d/40-resource-guardrails.conf
[Service]
MemoryHigh=300M
MemoryMax=400M
CPUQuota=90%
TasksMax=512
