# services/service_overrides/hermes-gateway.service.d/50-runtime-self-protection.conf
[Service]
Environment=HERMES_PROTECT_RUNTIME_COMMANDS=1
