# services/service_overrides/hermes-gateway.service.d/zz-runtime-first.conf
[Service]
WorkingDirectory=/home/Bilirubin/.hermes/hermes-agent
EnvironmentFile=-/home/Bilirubin/.hermes/.env
