# services/service_overrides/ralph.service.d/20-resource-guardrails.conf
[Service]
MemoryMax=256M
MemoryHigh=200M
CPUQuota=30%
TasksMax=50
