# services/service_overrides/infra-health-loop.timer.d/20-jitter.conf
[Timer]
RandomizedDelaySec=120
AccuracySec=1min
