# services/service_overrides/automation-gateway.service.d/20-resource-guardrails.conf
[Service]
MemoryMax=220M
CPUQuota=30%
TasksMax=256
