# services/hermes-self-monitor.service
[Unit]
Description=Hermes self-monitoring and remediation via MCP
After=network-online.target brain-mcp.service hermes-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
User=Bilirubin
Group=Bilirubin
WorkingDirectory=/home/Bilirubin/workspace
Environment=HOME=/home/Bilirubin
Environment=HERMES_HOME=/home/Bilirubin/.hermes
Environment=PATH=/home/Bilirubin/.hermes/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -lc '/home/Bilirubin/.hermes/venv/bin/hermes cron run --accept-hooks 63dbc9feeb3a && /home/Bilirubin/.hermes/venv/bin/hermes cron tick --accept-hooks'
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
MemoryMax=512M
MemoryHigh=400M
CPUQuota=60%
TasksMax=100
