# services/prometheus.service
[Unit]
Description=Prometheus Monitoring
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=Bilirubin
Group=Bilirubin
WorkingDirectory=/srv/automation
ExecStartPre=-/usr/bin/docker rm -f prometheus
ExecStart=/usr/bin/docker run --rm \
  --name prometheus \
  --network host \
  --memory 300m \
  --memory-reservation 200m \
  --cpus 0.30 \
  -v /home/Bilirubin/workspace/prometheus.yml:/etc/prometheus/prometheus.yml:ro \
  -v prometheus_data:/prometheus \
  prom/prometheus:latest \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/prometheus \
  --storage.tsdb.retention.time=15d \
  --web.listen-address=127.0.0.1:9090 \
  --web.enable-lifecycle
ExecStop=/usr/bin/docker stop prometheus
Restart=on-failure
RestartSec=10
MemoryMax=400M
CPUQuota=35%

[Install]
WantedBy=multi-user.target
