# services/service_overrides/gdrive-backup.timer.d/20-jitter.conf
[Timer]
RandomizedDelaySec=300
AccuracySec=1min
