# services/service_overrides/brain-mcp.service.d/20-resource-guardrails.conf
[Service]
MemoryMax=300M
CPUQuota=30%
TasksMax=256
