# services/service_overrides/hermes-monitor.service.d/20-resource-guardrails.conf
[Service]
MemoryMax=150M
CPUQuota=20%
TasksMax=128
