# services/service_overrides/hermes-monitor.service.d/30-report-policy.conf
[Service]
Environment=HERMES_MONITOR_SEND_HEALTHY_REPORTS=0
Environment=HERMES_MONITOR_SEND_STARTUP_REPORT=0
