# services/service_overrides/n8n.service.d/20-resource-guardrails.conf
[Service]
MemoryMax=650M
CPUQuota=70%
TasksMax=512
