# services/service_overrides/infra-snapshot.timer.d/20-jitter.conf
[Timer]
RandomizedDelaySec=90
AccuracySec=1min
