# services/service_overrides/hermes-self-monitor.timer.d/20-jitter.conf
[Timer]
RandomizedDelaySec=180
AccuracySec=1min
